# aws-dns/update-route53.timer
[Unit]
Description=Updates route53 on a timer.

[Timer]
OnCalendar=*:0/30
Unit=update-route53.service
